# native/keysharp-inputd/systemd/keysharp-inputd.service.in
[Unit]
Description=Keysharp privileged Linux input broker
Requires=keysharp-inputd.socket
After=keysharp-inputd.socket systemd-udevd.service

[Service]
Type=simple
ExecStart=@CMAKE_INSTALL_FULL_BINDIR@/keysharp-inputd --system-service
NoNewPrivileges=yes
ProtectClock=yes
ProtectControlGroups=yes
ProtectKernelLogs=yes
ProtectKernelModules=yes
ProtectKernelTunables=yes
ProtectSystem=strict
ReadWritePaths=/var/lib/keysharp-inputd /var/lib/keysharp-trust
StateDirectory=keysharp-inputd keysharp-trust
UMask=0077
DevicePolicy=closed
DeviceAllow=char-input r
DeviceAllow=/dev/uinput rw
