# native/keysharp-inputd/systemd/keysharp-inputd.socket
[Unit]
Description=Keysharp privileged Linux input broker socket

[Socket]
ListenStream=/run/keysharp-inputd/keysharp-inputd.sock
SocketMode=0666
DirectoryMode=0755
Accept=no

# The service requires this unit, so the installer enables only the service.
# Keep install metadata so upgrades/uninstalls can cleanly disable a socket that
# an older release enabled, without systemd diagnosing this unit as static.
[Install]
WantedBy=sockets.target
