# native/keysharp-inputd/systemd/keysharp-inputd.service.in
[Unit]
Description=Keysharp privileged Linux input broker
Requires=keysharp-inputd.socket
After=keysharp-inputd.socket systemd-udevd.service keyd.service
# Soft ordering hint only (no Wants=/Requires=): if keyd is installed and
# starts around the same time as this unit, let it grab raw physical devices
# first so we correctly end up downstream of it (grabbing keyd's re-emitted
# virtual keyboard, per protocol.h's vendor-masquerade design) rather than
# racing it for the raw device. Harmless no-op if keyd.service does not exist
# on this system -- systemd silently ignores an After= on a unit that isn't
# present. The socket remains a recovery activation path, so a manually-stopped
# daemon may still be started later by a client regardless of boot ordering.

[Service]
Type=simple
ExecStart=@CMAKE_INSTALL_FULL_BINDIR@/keysharp-inputd --system-service
# Recover promptly if the daemon crashes while devices are grabbed (the kernel
# releases grabs on exit, so this fails open; restart re-establishes hooking).
Restart=on-failure
RestartSec=1
# While hooks are active this daemon sits between the keyboard and the whole
# system; being paged out or OOM-killed then becomes system-wide input lag/loss.
# Deprioritize it for the OOM killer and permit its small, bounded resident set
# to be locked only for the lifetime of an interception grab.
OOMScoreAdjust=-900
LimitMEMLOCK=infinity
NoNewPrivileges=yes
ProtectClock=yes
ProtectControlGroups=yes
ProtectKernelLogs=yes
ProtectKernelModules=yes
ProtectKernelTunables=yes
ProtectSystem=strict
ReadWritePaths=/var/lib/keysharp-trust
StateDirectory=keysharp-trust
UMask=0077
DevicePolicy=closed
DeviceAllow=char-input r
DeviceAllow=/dev/uinput rw

[Install]
WantedBy=multi-user.target
